# packages/python-controller/systemd/vllm-llama.service
[Unit]
Description=vLLM Service for Llama Model
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/ai-os
ExecStart=/usr/local/bin/python -m vllm.entrypoints.api_server \
    --model /models/llama-3-8b-instruct \
    --port 8001 \
    --tensor-parallel-size 1 \
    --gpu-memory-utilization 0.85 \
    --max-model-len 8192 \
    --enable-chunked-prefill
Restart=always
RestartSec=10
Environment="PYTHONPATH=/root/ai-os"
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target